# packages/proxy/doc/pw-notifier/pm-notifier.service
[Unit]
Description=ProtonMail Notifier

[Service]
Type=simple
ExecStart=pm-notifier
Restart=on-failare

[Install]
WantedBy=default.target
